# local/service/syncthing/run
#!/bin/sh

exec syncthing --no-browser >/dev/null 2>&1
